# packages/bsp/common/etc/udev/rules.d/99-pad7-vibration.rules
# /sys/class/leds/pad7_vibration/brightness for user group
ACTION=="add", SUBSYSTEM=="leds", RUN+="/bin/chgrp video $sys$devpath/brightness", RUN+="/bin/chmod 0666 $sys$devpath/brightness"
